# files/supervisord.d/apache.conf
[program:httpd]
command=/etc/apache2/foreground.sh
stopsignal=6

