# files/supervisord.d/postfix.conf
[program:postfix]
command=/opt/postfix.sh
autostart=true
autorestart=true

